# systemd/.config/systemd/user/mbsync.timer
[Unit]
Description=Run mbsync and mu index every five minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
